# stow/systemd/.config/systemd/user/mousehunt-dailies.timer
[Unit]
Description=Auto send dailies

[Timer]
OnCalendar=12:0
Unit=mousehunt-dailies.service
Persistent=True

[Install]
WantedBy=timers.target
