# stow/systemd/.config/systemd/user/mousehunt-server.service
[Unit]
Description=Run server

[Service]
Type=simple
ExecStart=/home/angkeith/Documents/Web-Bot/scripts/server.sh

[Install]
WantedBy=default.target
